# magni_bringup/debian/magni-base.service.in
[Unit]
Requires=NetworkManager.service time-sync.target

[Service]
Type=simple
User=ubuntu
Environment=ROS_HOME=/home/ubuntu/.ros
ExecStart=/bin/sh -c '. /home/ubuntu/ros2_ws/install/setup.sh; . /etc/ubiquity/env.sh && env | sort && ros2 launch --wait -v magni_bringup base.launch'

[Install]
WantedBy=multi-user.target
