# magni_bringup/debian/magni-base.timer
[Unit]
Description="Run magni-base.service 1 minute after boot"

[Timer]
OnBootSec=1min
Unit=magni-base.service

[Install]
WantedBy=multi-user.target
